# services/crowd-247-runner.service
[Unit]
Description=IRIS Crowd Analytics 24/7 Runner
After=iris-backend.service
Requires=iris-backend.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Desktop/iris_command
ExecStart=/home/ubuntu/Desktop/iris_command/py_backend/.venv/bin/python -u /home/ubuntu/Desktop/iris_command/crowd_247_runner.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
